# Systemd/hello.service
[Unit]
Description=Ruby on Rails bank app.
After=network.target

[Service]
User=ubuntu
Group=ubuntu
Environment=RAILS_ENV=production
WorkingDirectory=/var/www/apps/hello/current
ExecStart=/bin/bash -lc 'bundle exec puma --config config/puma.rb >> log/production.log'
TimeoutSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
